# deploy/idm-scheduled-events.service
[Unit]
Description=TCS Identity — process due scheduled events (rename cutover, alias reminders/removal)
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/idm
ExecStart=/usr/bin/php /var/www/idm/bin/run_scheduled_events.php
TimeoutStartSec=900
Nice=10
